# sh/shadowsocks-server.service
[Unit]
Description=Shadowsocks Server
After=network.target

[Service]
ExecStart=$HOME/go/bin/shadowsocks-server -c $HOME/config.json
Restart=always

[Install]
WantedBy=multi-user.target
